# s6-update-init
### BEGIN INIT INFO
# Provides:          s6-update-init
# Required-Start:    $syslog
# Required-Stop:     $syslog
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: SmartSix Gateway autoupdate initscript
# Description:       This file is intended to auto start at boot time
#                    the auto update script for SmartSix Gateway applicatio
#                    it must e placed in /etc/init.d
### END INIT INFO

# Author: Fabrizio Guglielmino <guglielmino@gmail.com>
#


DESC="SmartSix Gateway Updater"
DAEMON=/opt/s6-gateway-app/update.sh
